# service/.config/systemd/user/rclone-pics.timer
[Unit]
Description=Daily sync of Pictures via rclone

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target

