# service/.config/systemd/user/rclone-docs.timer
[Unit]
Description=Daily sync of Documents via rclone

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target

