# services/.config/systemd/user/backup.service
[Unit]
Description=Backup Arch Linux

[Service]
Type=oneshot
ExecStart=%h/Scripts/backup.sh
