# services/.config/systemd/user/rclone-multi.timer
[Unit]
Description=Run rclone-multi bisync every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
