# services/.config/systemd/user/rclone-multi.service
[Unit]
Description=Rclone multi bisync job

[Service]
Type=oneshot
ExecStart=/home/conart/Scripts/rclone-bisync.sh
